# deployment/supervisord-goffee.conf
[program:goffee]
command=/start-goffee.sh
numprocs=1
autostart=true
autorestart=true
