# deployment/supervisord-tor.conf
[program:tor]
command=/start-tor.sh
numprocs=1
autostart=true
autorestart=true
